# config/crontab
*/5 * * * *  bash ./app/lib/tasks/weekly_export.sh $RAILS_ENV
